# raid.automount
[Automount]
Where=/mnt/md0

[Install]
WantedBy=multi-user.target
